# payload/systemd/llama-server.service
[Unit]
Description=PC-wide llama.cpp OpenAI-compatible server
After=network.target

[Service]
Type=simple
User=llama-cpp
Group=llama-cpp
ExecStart=/usr/local/bin/llama-manager serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
ReadWritePaths=/var/lib/llama.cpp/state /var/log/llama.cpp

[Install]
WantedBy=multi-user.target
