# payload/systemd/ubuntu-zombie-health.service
[Unit]
Description=Periodic Ubuntu Zombie health check

[Service]
Type=oneshot
User=__AGENT_USER__
Group=__AGENT_USER__
Environment=ZOMBIE_DIR=__ZOMBIE_DIR__
Environment=ZOMBIE_HEALTH_WARN_ONLY=1
ExecStart=__ZOMBIE_DIR__/bin/health-check
